# scoremaster-controller/packaging/scoremaster-controller.service
[Unit]
Description=ScoreMaster Controller Renderer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=scoreboard
Group=scoreboard
EnvironmentFile=/etc/scoremaster-controller/config.env
# Capabilities are attached to the binary via setcap in postinst
ExecStart=/usr/bin/scoremaster-controller ${SCOREMASTER_ARGS}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
